# inky.service

[Unit]
Description=Inky display service

[Service]
ExecStart=/home/pi/InkyHassDisplay/startdisplay.sh

[Install]
WantedBy=multi-user.target
